# systemd/dci-openstack-agent.timer
[Unit]
Description=Trigger the DCI Openstack Agent

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h 

[Install]
WantedBy=timers.target
